# supervisor/start.conf
[program:start]
command=sh /opt/bin/entry_point.sh & sh /home/seluser/cytoscape/start.sh & sh /home/seluser/noVNC/utils/launch.sh --vnc localhost:5900
